# files/etc/init.d/m9-rtr-agent
#!/bin/sh /etc/rc.common
# procd service for the M9 router dashboard-sync agent.
START=95
STOP=10
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /usr/local/sbin/m9-rtr-agent
    procd_set_param respawn 3600 5 0   # restart on crash: after 5s, give up after many
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

reload_service() {
    stop
    start
}
